# ansible/roles/nightstop/templates/codefabrik-nightstop.service.j2
[Unit]
Description=Code-Fabrik Nacht-Stopp

[Service]
Type=oneshot
ExecStart=/opt/codefabrik/scripts/nightstop.sh
