# ansible/roles/nightstop/templates/codefabrik-nightstart.service.j2
[Unit]
Description=Code-Fabrik Morgen-Start

[Service]
Type=oneshot
ExecStart=/opt/codefabrik/scripts/nightstart.sh
