# ansible/roles/nightstop/templates/codefabrik-nightstop.timer.j2
[Unit]
Description=Code-Fabrik Nacht-Stopp ({{ nightstop_hour }}:00)

[Timer]
OnCalendar=*-*-* {{ nightstop_hour }}:00:00
Persistent=true

[Install]
WantedBy=timers.target
